# projectroot/usr/lib/systemd/system/mnt-data.mount
[Unit]
Description=Mount data partition

[Mount]
What=/dev/disk/by-usage/data
Where=/mnt/data
