# local_src/usb-gadget/usb-ifupdown@.service
[Unit]
Description=USB network on %I

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/sbin/ifup %I
ExecStop=/sbin/ifdown %I
